# com.adamcagle.claw.vrt2.queue.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adamcagle.claw.vrt2.queue</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/adamcagle/CLAW/VRT2/jobs/scheduler_dispatch.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/adamcagle/CLAW/VRT2</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/adamcagle</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>TZ</key>
        <string>America/New_York</string>
    </dict>

    <key>StartCalendarInterval</key>
    <array>
        <!-- 2:00 AM ET = 11:00 PM PT -->
        <dict>
            <key>Hour</key>   <integer>23</integer>
            <key>Minute</key> <integer>0</integer>
        </dict>
        <!-- 5:00 AM ET = 2:00 AM PT -->
        <dict>
            <key>Hour</key>   <integer>2</integer>
            <key>Minute</key> <integer>0</integer>
        </dict>
        <!-- 6:00 AM ET = 3:00 AM PT -->
        <dict>
            <key>Hour</key>   <integer>3</integer>
            <key>Minute</key> <integer>0</integer>
        </dict>
        <!-- 6:30 AM ET = 3:30 AM PT -->
        <dict>
            <key>Hour</key>   <integer>3</integer>
            <key>Minute</key> <integer>30</integer>
        </dict>
        <!-- 9:31 AM ET = 6:31 AM PT -->
        <dict>
            <key>Hour</key>   <integer>6</integer>
            <key>Minute</key> <integer>31</integer>
        </dict>
        <!-- 10:00 AM ET = 7:00 AM PT -->
        <dict>
            <key>Hour</key>   <integer>7</integer>
            <key>Minute</key> <integer>0</integer>
        </dict>
        <!-- 12:00 PM ET = 9:00 AM PT -->
        <dict>
            <key>Hour</key>   <integer>9</integer>
            <key>Minute</key> <integer>0</integer>
        </dict>
        <!-- 2:00 PM ET = 11:00 AM PT -->
        <dict>
            <key>Hour</key>   <integer>11</integer>
            <key>Minute</key> <integer>0</integer>
        </dict>
        <!-- 4:05 PM ET = 1:05 PM PT -->
        <dict>
            <key>Hour</key>   <integer>13</integer>
            <key>Minute</key> <integer>5</integer>
        </dict>
        <!-- 4:30 PM ET = 1:30 PM PT -->
        <dict>
            <key>Hour</key>   <integer>13</integer>
            <key>Minute</key> <integer>30</integer>
        </dict>
        <!-- 6:00 PM ET = 3:00 PM PT -->
        <dict>
            <key>Hour</key>   <integer>15</integer>
            <key>Minute</key> <integer>0</integer>
        </dict>
        <!-- 9:15 PM ET = 6:15 PM PT -->
        <dict>
            <key>Hour</key>   <integer>18</integer>
            <key>Minute</key> <integer>15</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/adamcagle/CLAW/VRT2/logs/scheduler.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/adamcagle/CLAW/VRT2/logs/scheduler.stderr.log</string>
</dict>
</plist>
